# misc/freecloud/freecloud/deployment/supervisor.conf
[program:freecloud]
command=/usr/local/bin/gunicorn freecloud.wsgi:application -b 0.0.0.0:8000
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/freecloud.err.log
stdout_logfile=/var/log/freecloud.out.log

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/access.log
